# crontabs/root
# crontab
*/1 * * * * /usr/bin/watcher.sh >> /var/log/docker/watcher.log
*/5 * * * * /usr/bin/cleanup.sh >> /var/log/docker/cleanup.log
# 24 * * * * /usr/bin/buoy.sh >> /var/log/docker/buoy.log
42 1 * * * /usr/bin/vacuum.sh >> /var/log/docker/vacuum.log
*/4 * * * * /usr/bin/refresh.sh >> /var/log/docker/refresh.log
