# crontabs/root
# crontab
*/1 * * * * /usr/bin/watcher.sh >> /var/log/cron/watcher.log
*/5 * * * * /usr/bin/cleanup.sh >> /var/log/cron/cleanup.log
# 24 * * * * /usr/bin/buoy.sh >> /var/log/cron/buoy.log
42 1 * * * /usr/bin/vacuum.sh >> /var/log/cron/vacuum.log
*/4 * * * * /usr/bin/refresh.sh >> /var/log/cron/refresh.log
* * * * *   echo "Hello $(date)" >> /var/log/cron/cron.log 2>&1
